# systemd/ittsd.service
[Unit]
Description=iTTSd local queued streaming TTS daemon
After=default.target sound.target pipewire.service qwen3-fast-tts.service

[Service]
Type=simple
WorkingDirectory=%h/dev/ittsd
ExecStart=%h/dev/ittsd/bin/ittsd --addr 127.0.0.1:8765 --runtime-dir %h/.cache/ittsd --voice Vivian --lang auto --fast-url http://127.0.0.1:8001
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
